# hypraceelerator/systemd/clipboard-cleanup.service
[Unit]
Description=Hypraceelerator Clipboard Cleanup

[Service]
Type=oneshot
ExecStart=%h/.config/hypraceelerator/scripts/clipboard-cleanup.sh
