# hypraceelerator/systemd/clipboard-cleanup.timer
[Unit]
Description=Daily Clipboard Cleanup

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
