# hypraceelerator/systemd/dotfiles-backup.service
[Unit]
Description=Hypraceelerator Dotfiles Backup

[Service]
Type=oneshot
ExecStart=%h/.config/hypraceelerator/scripts/dotfiles-backup.sh
