# hypraceelerator/systemd/dotfiles-backup.timer
[Unit]
Description=Weekly Dotfiles Backup

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
