# hypraceelerator/systemd/downloads-archive.service
[Unit]
Description=Hypraceelerator Downloads Archive

[Service]
Type=oneshot
ExecStart=%h/.config/hypraceelerator/scripts/downloads-archive.sh
