# hypraceelerator/systemd/downloads-archive.timer
[Unit]
Description=Daily Downloads Archive

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
