# hypraceelerator/systemd/screenshots-organizer.timer
[Unit]
Description=Weekly Screenshot Organizer

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
