# hypraceelerator/systemd/system-alerts.service
[Unit]
Description=Hypraceelerator System Alerts

[Service]
Type=oneshot
ExecStart=%h/.config/hypraceelerator/scripts/system-alerts.sh
